# ss-tproxy.service
[Unit]
Description=ss/ssr/v2ray/socks5 transparent proxy script
ConditionFileIsExecutable=/usr/local/bin/ss-tproxy
ConditionFileNotEmpty=/etc/tproxy/ss-tproxy.conf
Requires=network.target network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/ss-tproxy start
ExecStop=/usr/local/bin/ss-tproxy stop
ExecReload=/usr/local/bin/ss-tproxy restart

[Install]
WantedBy=multi-user.target
